# deploy/supervisor/seoai-runtime.conf
[program:seoai-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/seoai/artisan queue:work --queue=crawl,generation,publishing,default --tries=3 --sleep=1 --timeout=120 --memory=512 --max-time=3600
directory=/var/www/seoai
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
startsecs=5
stopwaitsecs=360
stdout_logfile=/var/www/seoai/storage/logs/supervisor-queue.log
stderr_logfile=/var/www/seoai/storage/logs/supervisor-queue-error.log
redirect_stderr=false
environment=APP_ENV="production",QUEUE_CONNECTION="database"

[program:seoai-scheduler]
command=/usr/bin/php /var/www/seoai/artisan schedule:work
directory=/var/www/seoai
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=5
stopwaitsecs=120
stdout_logfile=/var/www/seoai/storage/logs/supervisor-scheduler.log
stderr_logfile=/var/www/seoai/storage/logs/supervisor-scheduler-error.log
redirect_stderr=false
environment=APP_ENV="production"